# rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv sh

set -ex
cat /etc/nginx/nginx.conf
exec nginx;
